# waypoint.service
[Unit]
Description=A simple caching and redirect backend

[Service]
Type=simple
ExecStart=/usr/bin/waypoint

[Install]
WantedBy=multi-user.target
